# inkforge.service
[Unit]
Description=InkForgeAPI - Chinese Handwriting Recognition API Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=owenx
Group=owenx
WorkingDirectory=/data/aozaiweb/InkForgeAPI

# 使用uvicorn启动ASGI应用
ExecStart=/data/aozaiweb/InkForgeAPI/venv/bin/uvicorn main:app --host 0.0.0.0 --port 20802 --workers 1

# 环境变量设置
Environment="PATH=/data/aozaiweb/InkForgeAPI/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH=/data/aozaiweb/InkForgeAPI"
Environment="OMP_NUM_THREADS=1"

# 重启策略
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=5

# 进程管理
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=30

# 资源限制（根据设备调整）
# MemoryMax=512M
# CPUQuota=50%

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=inkforgeapi

# 日志轮转配置（需要配合logrotate）
# 应用层日志文件在 /mnt/data/logs/inkforgeapi.log

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/data/aozaiweb/InkForgeAPI /tmp /data/.inkforge
ProtectHome=true

[Install]
WantedBy=multi-user.target